# openvpn_server/rootfs/etc/services.d/openvpn/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: OpenVPN Server
# Take down the S6 supervision tree when OpenVPN fails
# ==============================================================================

declare exit_code
readonly exit_code_container="${1}"

bashio::log.info "OpenVPN server stopped with exit code ${exit_code_container}"

if [[ "${exit_code_container}" -ne 0 ]] && [[ "${exit_code_container}" -ne 256 ]]; then
  bashio::log.warning "Halt add-on with exit code ${exit_code_container}"
  echo "${exit_code_container}" > /run/s6-linux-init-container-results/exitcode
  exec /run/s6/basedir/bin/halt
fi

bashio::log.info "Service restart after closing"
