# openvpn_server/rootfs/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: OpenVPN Server
# Secure PKI in /data, user files in /share
# ==============================================================================

# Private data in /data
declare data_dir="/data/openvpn"
declare pki_dir="${data_dir}/pki"
declare server_conf="${data_dir}/server.conf"

# User-accessible in /share
declare share_dir="/share/openvpn"
declare clients_dir="${share_dir}/clients"

bashio::log.info "Starting OpenVPN Server..."

# Create directory structure
mkdir -p "${data_dir}"
mkdir -p "${pki_dir}"
mkdir -p "${clients_dir}"

# Check if PKI is COMPLETELY initialized - verify ALL required files
if [[ ! -f "${pki_dir}/ca.crt" ]] || \
   [[ ! -f "${pki_dir}/issued/server.crt" ]] || \
   [[ ! -f "${pki_dir}/private/server.key" ]] || \
   [[ ! -f "${pki_dir}/dh.pem" ]] || \
   [[ ! -f "${pki_dir}/ta.key" ]]; then
    
    bashio::log.warning "PKI incomplete or missing. (Re)initializing..."
    
    # Clean up any incomplete PKI to start fresh
    rm -rf "${pki_dir}"
    rm -rf "${data_dir}/easy-rsa"
    
    # Run full initialization
    bashio::log.info "Running complete PKI initialization..."
    /usr/bin/openvpn-init || bashio::exit.nok "Failed to initialize PKI"
    
    # Verify initialization succeeded
    if [[ ! -f "${pki_dir}/ta.key" ]]; then
        bashio::exit.nok "PKI initialization failed - ta.key not created"
    fi
    
    bashio::log.info "✓ PKI initialization successful"
fi

# Get configuration
PROTO=$(bashio::config 'proto')
PORT=$(bashio::config 'port')
SERVER_SUBNET=$(bashio::config 'server_subnet')
SERVER_NETMASK=$(bashio::config 'server_netmask')
CIPHER=$(bashio::config 'cipher')
AUTH=$(bashio::config 'auth')
CLIENT_TO_CLIENT=$(bashio::config 'client_to_client')
DUPLICATE_CN=$(bashio::config 'duplicate_cn')
MAX_CLIENTS=$(bashio::config 'max_clients')
LOG_LEVEL=$(bashio::config 'log_level')

bashio::log.info "Generating server configuration..."

# Create server.conf
cat > "${server_conf}" <<EOF
# OpenVPN Server Configuration
# Generated by Home Assistant Add-on
port ${PORT}
proto ${PROTO}
dev tun0
topology subnet

# Certificates and keys (secure in /data)
ca ${pki_dir}/ca.crt
cert ${pki_dir}/issued/server.crt
key ${pki_dir}/private/server.key
dh ${pki_dir}/dh.pem
tls-auth ${pki_dir}/ta.key 0

# Network configuration
server ${SERVER_SUBNET} ${SERVER_NETMASK}
ifconfig-pool-persist ${data_dir}/ipp.txt

# Push routes to clients
EOF

# Add DNS servers
for dns in $(bashio::config 'dns_servers'); do
    echo "push \"dhcp-option DNS ${dns}\"" >> "${server_conf}"
done

# Add routes - iterate through array properly  
declare route_count
route_count=$(bashio::config 'routes | length')

if [[ "${route_count}" -gt 0 ]]; then
    for (( i=0; i < route_count; i++ )); do
        route=$(bashio::config "routes[${i}]")
        echo "push \"route ${route}\"" >> "${server_conf}"
    done
fi

# Client options
if bashio::config.true 'client_to_client'; then
    echo "client-to-client" >> "${server_conf}"
fi

if bashio::config.true 'duplicate_cn'; then
    echo "duplicate-cn" >> "${server_conf}"
fi

# Continue server.conf
cat >> "${server_conf}" <<EOF

# Connection settings
max-clients ${MAX_CLIENTS}
keepalive 10 120
persist-key
persist-tun

# Security
cipher ${CIPHER}
auth ${AUTH}
tls-version-min 1.2

# Logging
status ${data_dir}/openvpn-status.log
verb ${LOG_LEVEL}

# Client config directory
client-config-dir ${data_dir}/ccd

# Explicit exit notify for UDP
explicit-exit-notify 1
EOF

# Create client config directory
mkdir -p "${data_dir}/ccd"

# Create user README in /share
cat > "${share_dir}/README.txt" <<'READMEEOF'
OpenVPN Server - Download Your Clients Here
============================================

Client configuration files (.ovpn) are in the "clients" folder.

Download via:
- File Editor: Browse to /share/openvpn/clients/
- Samba: \\homeassistant\share\openvpn\clients\
- SSH: /share/openvpn/clients/

Before using a client file:
1. Open the .ovpn file in a text editor
2. Find the line: remote YOUR_PUBLIC_IP_OR_DDNS 1194
3. Replace YOUR_PUBLIC_IP_OR_DDNS with your actual:
   - Public IP address (find at whatismyip.com), OR
   - DuckDNS hostname (e.g., yourhome.duckdns.org)
4. Save the file
5. Import to OpenVPN Connect app on your device

To create new clients:
1. Open OpenVPN Server add-on
2. Click "Terminal" tab (appears once add-on is running)
3. Run: openvpn-client-add clientname
4. Download the .ovpn file from clients folder

Commands available in Terminal:
- openvpn-client-add <name>     Create new client
- openvpn-client-list            List all clients
- openvpn-client-revoke <name>   Revoke a client

Port forwarding required:
Forward UDP port 1194 on your router to your Home Assistant IP.
READMEEOF

# Enable IP forwarding
bashio::log.info "Enabling IP forwarding..."
sysctl -w net.ipv4.ip_forward=1 >/dev/null 2>&1 || \
    bashio::log.warning "Could not enable IP forwarding (may already be enabled by system)"

# Setup iptables rules
bashio::log.info "Configuring firewall rules..."
iptables -t nat -A POSTROUTING -s ${SERVER_SUBNET}/24 -o eth0 -j MASQUERADE 2>&1 || true
iptables -A FORWARD -i tun0 -j ACCEPT 2>&1 || true
iptables -A FORWARD -o tun0 -j ACCEPT 2>&1 || true

bashio::log.info "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
bashio::log.info "✓ PKI secure in: /data/openvpn/pki/"
bashio::log.info "✓ Download clients from: /share/openvpn/clients/"
bashio::log.info "✓ Server config: ${PROTO}/${PORT}"
bashio::log.info "✓ VPN subnet: ${SERVER_SUBNET}/${SERVER_NETMASK}"
bashio::log.info "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
bashio::log.info "Starting OpenVPN daemon..."

# Check if client creation requested
if bashio::config.true 'create_client_button' && [[ -n "$(bashio::config 'client_name')" ]]; then
    bashio::log.info "Client creation requested via config: $(bashio::config 'client_name')"
    /usr/bin/openvpn-client-add "$(bashio::config 'client_name')"
    bashio::log.info "Client creation completed. Restart add-on to use normally."
    # Don't start OpenVPN - let user restart after config change
    bashio::exit.ok "Client created successfully"
fi

# Start OpenVPN - redirect stderr to stdout to capture errors
openvpn --config "${server_conf}" 2>&1
