# openvpn_server/rootfs/etc/services.d/health/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
exec /usr/bin/health-server
