# openvpn_server/rootfs/etc/services.d/sensor-update/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
sleep 15
exec /usr/bin/openvpn-sensor-update
